# mqtt.cool/bin/unix-like/install/OpenBSD/init/mqttcool
#!/bin/sh
#
# $OpenBSD: $

daemon="%MC_HOME%/bin/unix-like/mc.sh"
daemon_flags=""
daemon_user=%MC_USER%
daemon_pid="%MC_HOME%/bin/unix-like/daemon.pid"
# RC_DEBUG=1

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
    ${rcexec} "WRITE_PID=${daemon_pid} ${daemon} run &"
}

rc_check() {
    _pid=$(cat ${daemon_pid})
    pgrep -u "${daemon_user}" | grep "^${_pid}" > /dev/null
}

rc_stop() {
    _pid=$(cat ${daemon_pid})
    kill -TERM "${_pid}"
}

rc_cmd $1
